# hp_monasca_ansible/roles/storm/files/storm-logviewer.conf
# Startup script for Storm Logviewer

description "Storm Logviewer daemon"
start on runlevel [2345]

console log
respawn

kill timeout 240
respawn limit 25 5

setgid storm
setuid storm
chdir /opt/storm/current
exec /opt/storm/current/bin/storm logviewer
