# hp_ceph_ansible/roles/_CEP-CMN/files/ceph.target
[Unit]
Description=ceph target allows starting and stopping all ceph*@.service instances at once
[Install]
WantedBy=multi-user.target
